# sensor/systemd/jett.service
[Unit]
Description=jeTT Real-Time eBPF AI Verdict Daemon
After=kage-sensor.service bifrost-guardian.service
Requires=kage-sensor.service bifrost-guardian.service

[Service]
Type=simple
ExecStart=/usr/bin/jett_daemon
Restart=always
RestartSec=1
User=root
Group=root

[Install]
WantedBy=multi-user.target
